# com.bee.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bee.sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-l</string>
        <string>-c</string>
        <string>/Users/YOUR_USER/bee-sync/sync.sh</string>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USER/bee-sync/sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USER/bee-sync/sync.log</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
</dict>
</plist>
